# deploy/systemd/velocity-backup.timer
[Unit]
Description=Run Velocity backup nightly at 04:30

[Timer]
OnCalendar=*-*-* 04:30:00
Persistent=true

[Install]
WantedBy=timers.target
