# deploy/systemd/velocity-backup.service
[Unit]
Description=Velocity nightly pg_dump

[Service]
Type=oneshot
User=yusuf
ExecStart=/home/yusuf/velocity/deploy/backup/pg_backup.sh
