# deploy/systemd/velocity.service
[Unit]
Description=Velocity game server (FastAPI/uvicorn)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=yusuf
WorkingDirectory=/home/yusuf/velocity/server
Environment=PYTHONUNBUFFERED=1
ExecStartPre=/usr/bin/env bash -c 'cd /home/yusuf/velocity/deploy && docker compose up -d'
ExecStartPre=/home/yusuf/velocity/server/.venv/bin/alembic upgrade head
ExecStart=/home/yusuf/velocity/server/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
